# server/titanic.service
[Unit]
Description=Titanic dj website

[Service]
ExecStart=/home/ubuntu/server.sh

[Install]
WantedBy=multi-user.target